# deploy/gunicorn-excavator.service
[Unit]
Description=gunicorn daemon for excavator
After=network.target

[Service]
User=root
Group=www-data
WorkingDirectory=/srv/excavator
Environment="PATH=/srv/excavator/venv/bin"
ExecStart=/srv/excavator/venv/bin/gunicorn config.wsgi:application --workers 3 --timeout 60 --bind unix:/run/gunicorn-excavator.sock --access-logfile - --error-logfile -
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
